# autostart/1password.desktop
[Desktop Entry]
Comment=Manage your passwords
Exec=1password
GenericName=Password Manager
Hidden=true
Name=1Password
Type=Application
