# autostart/1password.desktop
[Desktop Entry]
Type=Application
Name=1Password
GenericName=Password Manager
Comment=Manage your passwords
Exec=1password
